# shell/ubuntu/16/axe-aggregator.service
[Unit]
Description=Axe-Aggregator


[Service]
Type=simple
ExecStart=/opt/cloudiator/start_axe-aggregator.sh

[Install]
WantedBy=multi-user.target

