# shell/ubuntu/16/colosseum.service
[Unit]
Description=Colosseum


[Service]
Type=simple
ExecStart=/opt/cloudiator/start_colosseum.sh

[Install]
WantedBy=multi-user.target



