# shell/ubuntu/16/etcd.service
[Unit]
Description=etcd


[Service]
Type=simple
ExecStart=/opt/cloudiator/start_etcd.sh

[Install]
WantedBy=multi-user.target


